# chatdku-deploy.timer
[Unit]
Description=Periodically check if ChatDKU frontend needs redeployment

[Timer]
OnBootSec=2min
OnUnitActiveSec=15min

[Install]
WantedBy=timers.target
