# chatdku-deploy.service
[Unit]
Description=Deploy ChatDKU frontend if missing from /var/www/chatdku
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/home/anar/Developer/ChatDKU-web/chatdku-deploy.sh
User=anar
Environment=HOME=/home/anar
Environment=PATH=/usr/local/bin:/usr/bin:/bin

[Install]
WantedBy=multi-user.target
